# mac/claude-wake.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.__USERNAME__.claude-wake</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>__HOME__/.claude-wake-launcher.sh</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>__HOME__/.claude-wake-trigger</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
